# docker/dovecot-postfix/build/supervisord/conf.d/dovecot.conf
[program:dovecot]
command=/usr/local/sbin/dovecot -F
priority=900
autostart=true
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
