# docker/dovecot-postfix/build/supervisord/conf.d/postfix.conf
[program:postfix]
command=/usr/local/postfix/sbin/postfix start
priority=900
autostart=true
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
